# hermes-webui.service
[Unit]
Description=Hermes Web UI Backend
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/house/vibecode
ExecStart=/usr/bin/python3 /home/house/vibecode/server.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
