# vibecode.service
[Unit]
Description=vibecode server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# User=house (optional; set to your runtime user if needed)
WorkingDirectory=%h/vibecode
ExecStart=/usr/bin/python3 %h/vibecode/server.py
Restart=always
RestartSec=5
StartLimitInterval=60
StartLimitBurst=3
Environment=HOME=%h
Environment=HERMES_WEBUI_PORT=8786
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONOPTIMIZE=1
Environment=PYTHONPATH=%h/vibecode
StandardOutput=journal
StandardError=journal
# Performance optimizations
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=default.target
