# tools/auto_push.service
[Unit]
description=vibecode auto-sync service
after=network-online.target
wants=network-online.target

[Service]
type=simple
workingdirectory=%h/vibecode
execstart=%h/vibecode/tools/auto_push.sh
restart=always
restartsec=10
standardoutput=journal
standarderror=journal

environment=GIT_SYNC_BRANCH=main

[Install]
wantedby=default.target
